# docker-support/docker-service
#!/sbin/openrc-run
# Docker service for NASBox

name="Docker"
description="Docker Application Container Engine"

command="/usr/bin/dockerd"
command_args="--config-file /etc/docker/daemon.json"
command_background="yes"
pidfile="/var/run/docker.pid"
start_stop_daemon_args="--background --make-pidfile"

depend() {
    need net
    after firewall
}

start_pre() {
    # Ensure required directories exist
    checkpath --directory --mode 0755 /var/lib/docker
    checkpath --directory --mode 0755 /var/run/docker
    
    # Load required kernel modules
    modprobe overlay
    modprobe br_netfilter
    
    # Enable IP forwarding
    sysctl -w net.ipv4.ip_forward=1 > /dev/null 2>&1
    sysctl -w net.bridge.bridge-nf-call-iptables=1 > /dev/null 2>&1
}

healthcheck() {
    docker info > /dev/null 2>&1
}
